# home-dashboard/rootfs/etc/s6-overlay/s6-rc.d/dashboard/run
#!/command/with-contenv bashio
set -e

export HA_TOKEN=$(bashio::config 'ha_token')
export HA_BASE_URL="http://homeassistant:8123"
export PORT=8099
export HOSTNAME=0.0.0.0

cd /app
exec node server.js
